# configs/code-server.service
[Unit]
Description=code-server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i
# Enumerate the user's tool dirs explicitly so code-server's integrated
# terminal and tasks can find npm globals, cargo, bun, uv, go, etc.
# (We can't rely on `bash -lc` here — bash login shells don't source zshrc.)
Environment=PATH=/home/%i/.local/bin:/home/%i/.npm-global/bin:/home/%i/.bun/bin:/home/%i/.cargo/bin:/home/%i/go/bin:/usr/local/bin:/usr/bin
ExecStart=/usr/bin/code-server
Restart=always

[Install]
WantedBy=multi-user.target
